# monitoring/node_exporter/node_exporter.service
# Runs on each Proxmox host (pve, node2, node3) - not in a container
# Reads system stats from /proc and /sys, exposes them on port 9100
# Install: cp to /etc/systemd/system/, then systemctl daemon-reload && systemctl enable --now node_exporter

[Unit]
Description=Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/node_exporter
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
