# monitoring/systemd/alertmanager-ntfy.service
# alertmanager-ntfy.service - Runs inside the monitoring CT
#
# Translates Alertmanager webhook JSON into ntfy push notifications
# Without this bridge, Alertmanager can't talk to ntfy directly
#
# Install location: /opt/alertmanager-ntfy/alertmanager-ntfy
# Config: /opt/alertmanager-ntfy/config.yml
# Listens on: port 8086

[Unit]
Description=Alertmanager to ntfy bridge
Documentation=https://github.com/alexbakker/alertmanager-ntfy
After=network.target

[Service]
Type=simple

ExecStart=/opt/alertmanager-ntfy/alertmanager-ntfy \
  --config /opt/alertmanager-ntfy/config.yml

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
