# monitoring/systemd/alertmanager.service
# Runs in the monitoring LXC Container
# Receives alerts from Prometheus, groups them, sends to PagerDuty
# Web UI: http://<CT-IP>:9093

[Unit]
Description=Alertmanager
Documentation=https://prometheus.io/docs/alerting/alertmanager/
After=network.target

[Service]
Type=simple
ExecStart=/opt/alertmanager/alertmanager \
  --config.file=/opt/alertmanager/alertmanager.yml \
  --storage.path=/opt/alertmanager/data
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
