# monitoring/systemd/blackbox_exporter.service
# Runs in the monitoring CT (110)
# Performs HTTP, TCP, and ICMP checks when Prometheus asks it to
# Web UI: http://<CT-IP>:9115

[Unit]
Description=Blackbox Exporter
Documentation=https://github.com/prometheus/blackbox_exporter
After=network.target

[Service]
Type=simple
ExecStart=/opt/blackbox_exporter/blackbox_exporter \
  --config.file=/opt/blackbox_exporter/blackbox.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
