# monitoring/systemd/prometheus.service
# Runs in the monitoring CT (110)
# Scrapes all targets, stores metrics for 90 days, evaluates alert rules
# Web UI: http://<CT-IP>:9090

[Unit]
Description=Prometheus Monitoring System
Documentation=https://prometheus.io/docs/
After=network.target

[Service]
Type=simple
ExecStart=/opt/prometheus/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/opt/prometheus/data \
  --storage.tsdb.retention.time=90d \
  --web.enable-lifecycle
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
